# rclone-manager.service
[Unit]
Description=RClone Manager Web Application
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/Rclone-Manager

# Sostituzione di ExecStart per eseguire start_all.sh con sorgente venv
ExecStart=/bin/bash -c 'source /opt/Rclone-Manager/venv/bin/activate && ./start_all.sh'

# Make sure directories exist
ExecStartPre=/bin/mkdir -p /opt/Rclone-Manager/data/logs
ExecStartPre=/bin/chmod +x /opt/Rclone-Manager/start_all.sh

# Imposta variabili d'ambiente
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_ENV=production

# Riavvia il servizio in caso di errore
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# Configurazione logging
SyslogIdentifier=rclone-manager
StandardOutput=journal
StandardError=journal

# Limiti di sicurezza e isolamento
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
